# etc/chaos_supervisor.conf
[program:chaos]
command=/root/.virtualenvs/chaos/bin/python chaos.py
directory=/root/workspace/Chaos
user=root
